# deploy/crontab
* *   * * *   www-ticketing   /var/www/2013/live/venv/bin/python /var/www/2013/live/bin/live_receipt.py
